# DEBIAN/postinst
#!/bin/sh
systemctl enable ossec-metrics
systemctl start ossec-metrics
