# DEBIAN/prerm
#!/bin/sh
systemctl stop ossec-metrics
systemctl disable ossec-metrics
